# grott/rootfs/etc/services.d/grott/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Grott
# Runs Grott
# ==============================================================================
bashio::log.info "Starting Grott..."
exec /usr/bin/grott